# ops/systemd/gestiq-worker-imports.service
[Unit]
Description=GestiQ Imports Celery Worker
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=gestiq
Group=gestiq
WorkingDirectory=/opt/gestiq/apps/backend
Environment="PYTHONPATH=/opt/gestiq/apps/backend:/opt/gestiq/apps"
Environment="DB_DSN=postgresql://gestiq:PASSWORD@localhost:5432/gestiqcloud"
Environment="REDIS_URL=redis://localhost:6379/0"
Environment="IMPORTS_ENABLED=1"
Environment="IMPORTS_RUNNER_MODE=celery"
ExecStart=/opt/gestiq/venv/bin/python -m app.modules.imports.application.worker_main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gestiq-imports-worker

[Install]
WantedBy=multi-user.target
